# docker/supervisord.conf
; Supervisor configuration for ProtonVPN + Tailscale Exit Node
; Manages multiple processes within the container

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

; WireGuard process
[program:wireguard]
command=/app/wireguard-setup.sh setup
autostart=true
autorestart=false
startsecs=0
priority=1
stdout_logfile=/var/log/supervisor/wireguard.log
stderr_logfile=/var/log/supervisor/wireguard-error.log

; Tailscale daemon
[program:tailscaled]
command=/usr/local/bin/tailscaled --state=/var/lib/tailscale/tailscaled.state --socket=/var/run/tailscale/tailscaled.sock
autostart=true
autorestart=true
startsecs=5
priority=2
stdout_logfile=/var/log/supervisor/tailscaled.log
stderr_logfile=/var/log/supervisor/tailscaled-error.log

; Tailscale up command (runs once after tailscaled starts)
[program:tailscale-up]
command=/bin/bash -c 'sleep 5 && tailscale up --authkey=%(ENV_TAILSCALE_AUTH_KEY)s --hostname=%(ENV_TAILSCALE_HOSTNAME)s --advertise-exit-node --accept-dns=%(ENV_TAILSCALE_ACCEPT_DNS)s'
autostart=true
autorestart=false
startsecs=10
priority=3
stdout_logfile=/var/log/supervisor/tailscale-up.log
stderr_logfile=/var/log/supervisor/tailscale-up-error.log

; SOCKS5 Proxy (Dante)
[program:socks5-proxy]
command=/usr/sbin/sockd -f /etc/dante.conf -N 4
autostart=true
autorestart=true
startsecs=5
priority=4
stdout_logfile=/var/log/supervisor/socks5.log
stderr_logfile=/var/log/supervisor/socks5-error.log

; HTTP Proxy (TinyProxy)
[program:http-proxy]
command=/usr/bin/tinyproxy -d -c /etc/tinyproxy/tinyproxy.conf
autostart=true
autorestart=true
startsecs=5
priority=4
stdout_logfile=/var/log/supervisor/tinyproxy.log
stderr_logfile=/var/log/supervisor/tinyproxy-error.log

; Health check monitor
[program:health-monitor]
command=/bin/bash -c 'while true; do /app/entrypoint.sh healthcheck || echo "Health check failed at $(date)"; sleep 30; done'
autostart=true
autorestart=true
startsecs=60
priority=5
stdout_logfile=/var/log/supervisor/health-monitor.log
stderr_logfile=/var/log/supervisor/health-monitor-error.log
